# sample/collector.service
[Unit]
Description=collector service

[Service]
Type=simple
User=root
ExecStart=/opt/collector/collector

[Install]
WantedBy=multi-user.target